# scripts/cn.yiming1234.am-i-ok.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>cn.yiming1234.am-i-ok</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/path/to/am-i-ok-agent.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AM_I_OK_SECRET</key>
    <string>your_secret_here</string>
    <key>BETTER_AUTH_URL</key>
    <string>https://yiming1234.cn</string>
  </dict>

  <key>StartInterval</key>
  <integer>30</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardErrorPath</key>
  <string>/tmp/am-i-ok.log</string>
</dict>
</plist>